# laptalk.service
[Unit]
Description=LapTalk - Real-time voice transcription
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/make -C %h/src/laptalk run key=alt_r log=%h/src/laptalk/laptalk-log.txt
Restart=on-failure
RestartSec=5

# X11 environment (may need adjustment for your system)
# Check with: echo $DISPLAY and echo $XAUTHORITY
# Common values: DISPLAY=:0 or :1, XAUTHORITY varies by desktop environment
Environment="DISPLAY=:1"
Environment="XAUTHORITY=/run/user/1000/gdm/Xauthority"

[Install]
WantedBy=default.target
